# install/conf_files/ramdisk.service
[Unit]
Description=Ramdisk
After=syslog.target network.target

[Service]
ExecStart=/opt/integralstor/integralstor_common/install/scripts/ramdisk start 
User=root
TasksMax=infinity

[Install]
WantedBy=multi-user.target

