# install/conf_files/zed.conf
# zed - the ZFS  event daemon

description "The ZFS event daemon"

start on local-filesystems
stop on runlevel [!2345]

respawn
expect daemon

exec /sbin/zed
